# teamd/redhat/systemd/teamd@.service
[Unit]
Description=Team Daemon for device %I

[Service]
BusName=org.libteam.teamd.%i
ExecStart=/usr/bin/teamd -U -D -t %i -f /run/teamd/%i.conf
